# linux/systemd/eye-rest.timer
[Unit]
Description=Eye rest reminder

[Timer]
# Run every 15 minutes
OnCalendar=*:0/15
Persistent=false

[Install]
WantedBy=timers.target
