# linux/runit/lenovo_fix/run
#!/bin/sh

exec 2>&1
exec /usr/lib/throttled/lenovo_fix.py
#vim:ft=sh
